# systemd/system/snap-snapd-19457.mount
[Unit]
Description=Mount unit for snapd, revision 19457
After=snapd.mounts-pre.target
Before=snapd.mounts.target
Before=local-fs.target

[Mount]
What=/var/lib/snapd/snaps/snapd_19457.snap
Where=/snap/snapd/19457
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
